# files/geoipupdate.cron
#!/usr/bin/env bash
MAILTO=root
geoipupdate > /dev/null